# filedb-cleanup.timer
[Unit]
Description=Daily removal of unused files in filedb

[Timer]
OnCalendar=daily
RandomizedDelaySec=5h

[Install]
WantedBy=timers.target
