# files/filedb.service
[Unit]
Description=File DB

[Service]
User=filedb
Group=filedb
ExecStart=/usr/local/bin/filedbd

[Install]
WantedBy=wsgi-services.target
